# rpi_rootfs/etc/inittab
::sysinit:/etc/init.d/rcS

#user defined banner 
::once:/usr/bin/banner

console::askfirst:/bin/ash


#do this before poweroff
::ctrlaltdel:/sbin/poweroff
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
